# etc/systemd/system/systemd-timesyncd.service.d/override.conf
# NOTE: fix clock sync bug of WSL
# https://github.com/microsoft/WSL/issues/8204#issuecomment-1339506778
[Unit]
ConditionVirtualization=
