# src/dot_config/systemd/user/reload-waybar.path
[Unit]
Description=Watch waybar for changes

[Path]
PathModified=%h/.config/waybar/
PathModified=%h/.config/waybar/lib/

[Install]
WantedBy=multi-user.target
