# home/dot_config/systemd/user/reload-swaync.path
[Unit]
Description=Watch swaync for changes

[Path]
PathModified=%h/.config/swaync/

[Install]
WantedBy=graphical-session.target
