# home/dot_config/systemd/user/reload-kanshi.path
[Unit]
Description=Watch kanshi for changes

[Path]
PathModified=%h/.config/kanshi/config

[Install]
WantedBy=graphical-session.target
